# supervisord.example.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:octane]
command=php /code/artisan octane:start --server=frankenphp --host=0.0.0.0 --port=9000 --workers=2 --max-requests=250
directory=/code
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/octane.log
stderr_logfile=/var/log/supervisor/octane-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:reverb]
command=php /code/artisan reverb:start --host=0.0.0.0 --port=8006
directory=/code
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/reverb.log
stderr_logfile=/var/log/supervisor/reverb-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /code/artisan queue:work --queue=exports,default --sleep=3 --tries=3 --max-time=3600
directory=/code
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-worker.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=3600
killasgroup=true
stopasgroup=true

[program:scheduler]
command=php /code/artisan schedule:work
directory=/code
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=10MB

[program:laravel-bootstrap]
command=bash -c "php artisan config:clear && php artisan package:discover --ansi && php artisan config:cache && php artisan route:cache"
directory=/code
autostart=true
autorestart=false
startsecs=0
priority=1
stdout_logfile=/var/log/supervisor/bootstrap.log
stderr_logfile=/var/log/supervisor/bootstrap-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
